# .config/systemd/user/personaldocs.timer
[Unit]
Description=Sync Personal Documents daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
