# .config/systemd/user/university.timer
[Unit]
Description=Sync University Stuff daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
